# config/crontab
# m h dom mon dow user	command
15  */4   *  *  *   root  cron_task dump_postgres
